# config/docker/files/supervisord-dev.conf
[supervisord]
nodaemon=true

[program: solr]
command=bundle exec rake solr:run
umask=022

[program: unicorn]
command=bundle exec unicorn -c /app/config/docker/files/unicorn.rb
umask=022

