# config/docker/files/supervisord.conf
[supervisord]
nodaemon=true

[program: nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
umask=022

[program: unicorn]
command=/bin/bash -c "ln -s /app /var/www/app && cd /var/www/app && bundle exec unicorn -c /var/www/app/config/docker/files/unicorn.rb"
umask=022

